# dollar.conf
[program: dollar]
directory=/home/dollar
command=python app.py

autostart=true
autorestart=true
stderr_logfile=/home/dollar/logfile_err.log
stdout_logfile=/home/dollar/logfile.log
